# feishu-bridge.service
[Unit]
Description=Feishu Bridge
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
Type=simple
User=feishu-bridge
Group=feishu-bridge
WorkingDirectory=/opt/feishu-bridge
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/etc/feishu-bridge.env
ExecStart=/usr/bin/node /opt/feishu-bridge/bridge.js
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
